# podman/eyetor-telegram.container
# Quadlet Container unit — Eyetor Telegram bot (24/7 service)
# Place in: ~/.config/containers/systemd/
#
# After placing files, reload and enable:
#   systemctl --user daemon-reload
#   systemctl --user enable --now eyetor-telegram.service
#
# View logs:
#   journalctl --user -u eyetor-telegram -f

[Unit]
Description=Eyetor Telegram Bot
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/eyetor:latest
ContainerName=eyetor-telegram

# Environment — reads from ~/.config/eyetor/.env
EnvironmentFile=%h/.config/eyetor/.env

# Persistent data volume
Volume=eyetor-data.volume:/home/eyetor/.eyetor
# Skills and config from project directory (adjust path as needed)
Volume=%h/eyetor/skills:/app/skills:ro,z
Volume=%h/eyetor/config:/app/config:ro,z

# ---- HOST ACCESS (optional) ----
# Uncomment to give the agent access to specific directories.
# On Linux/macOS: use standard paths
# On Windows (WSL2): /mnt/c/Users/... maps to C:\Users\...
#
# Volume=%h/Documents:/workspace/documents:rw,z
# Volume=%h/Dev:/workspace/dev:rw,z

# Run the Telegram bot
Exec=eyetor telegram

# Health check
HealthCmd=eyetor providers list
HealthInterval=60s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=15s

# Auto-restart on failure
AutoUpdate=registry

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=60

[Install]
WantedBy=default.target
