# podman/eyetor-agent.container
# Quadlet Container unit — Eyetor Agent daemon (headless, 24/7)
# Use this when you want to expose the agent via a persistent process
# without an interactive channel (e.g. to connect from other systems).
#
# Place in: ~/.config/containers/systemd/
# systemctl --user enable --now eyetor-agent.service

[Unit]
Description=Eyetor Agent Daemon
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/eyetor:latest
ContainerName=eyetor-agent

EnvironmentFile=%h/.config/eyetor/.env

Volume=eyetor-data.volume:/home/eyetor/.eyetor
Volume=%h/eyetor/skills:/app/skills:ro,z
Volume=%h/eyetor/config:/app/config:ro,z

# Run a simple keep-alive: providers health check loop
# Override this with your own daemon command (e.g. a FastAPI server)
Exec=sh -c "while true; do eyetor providers list > /dev/null 2>&1; sleep 30; done"

HealthCmd=eyetor providers list
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
